# config/autostart/ibus-daemon.desktop
[Desktop Entry]
Type=Application
Name=IBus Daemon
Exec=ibus-daemon -drx
X-GNOME-Autostart-enabled=true